# templates/etc/profile.d/spark.sh.j2
export PATH=$PATH:{{spark_parent_dir}}/spark/bin
